# web-survey/deploy/launchd/com.berrie.survey.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.berrie.survey</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/ABSOLUTE/PATH/TO/web-survey/scripts/start_backend.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/ABSOLUTE/PATH/TO/web-survey</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>SURVEY_PORT</key>
    <string>8000</string>
    <key>SURVEY_DB_PATH</key>
    <string>/ABSOLUTE/PATH/TO/web-survey/survey_data.db</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/berrie-survey.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/berrie-survey.err.log</string>
</dict>
</plist>
